# config/systemd/user/hetzner-backup.service
[Unit]
Description=Hetzner Backup

[Service]
Type=oneshot
ExecStart=/home/mbezjak/workspace/dotfiles/bin/hetzner-backup.sh

[Install]
WantedBy=multi-user.target
